# filesystem/udev/rules/40-modprobe.rules
# hotplug
ENV{MODALIAS}=="?*",	RUN+="/sbin/modprobe $env{MODALIAS}"
